# root/etc/systemd/system/growroot.service
[Unit]
Description=Grow root partition
After=-.mount

[Service]
ExecStart=/usr/local/sbin/growroot

[Install]
WantedBy=default.target
